# scripts/memory-sync.timer
[Unit]
Description=Run Memory Sync every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
